# runit/monit/run
#!/bin/bash -le
exec 2>&1

exec /usr/bin/monit -c /etc/monitrc